# mtproto-proxy.service
[Unit]
Description=MTProto Proxy Server with DPI Bypass
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/aiagent/mtproto-proxy
ExecStart=/usr/bin/python3 -u /home/aiagent/mtproto-proxy/mtproto_proxy.py
Restart=always
RestartSec=5
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtproto-proxy

[Install]
WantedBy=multi-user.target
